# root/etc/services.d/TwitchPlaysX/run
#!/usr/bin/with-contenv sh

cd /opt/TwitchPlaysX && npm run start
